# sirvasms/supervisord/receiver.ini
[program:receiver]
command=python receiver.py
directory=/opt/djangoprojects/sirvasms/sirvasmsapp/scripts/
autostart=true
autorestart=true
stderr_logfile=/var/log/receiver.err.log
stdout_logfile=/var/log/receiver.out.log
