# sirvasms/supervisord/sender_special1.ini
[program:sender_special1]
command=python sender_special1.py
directory=/opt/djangoprojects/sirvasms/sirvasmsapp/scripts/
autostart=true
autorestart=true
stderr_logfile=/var/log/sender_special1.err.log
stdout_logfile=/var/log/sender_special1.out.log
