# sirvasms/supervisord/sender.ini
[program:sender]
command=python sender.py
directory=/opt/djangoprojects/sirvasms/sirvasmsapp/scripts/
autostart=true
autorestart=true
stderr_logfile=/var/log/sender.err.log
stdout_logfile=/var/log/sender.out.log
