# sirvasms/supervisord/sender_special2.ini
[program:sender_special2]
command=python sender_special2.py
directory=/opt/djangoprojects/sirvasms/sirvasmsapp/scripts/
autostart=true
autorestart=true
stderr_logfile=/var/log/sender_special2.err.log
stdout_logfile=/var/log/sender_special2.out.log
